# tablet-mode-accel-chuwi.service
[Unit]
Description=Software tablet mode switch
Before=gdm.service

[Service]
Type=exec
ExecStart=/usr/local/bin/tablet-mode-accel --chuwi

[Install]
WantedBy=graphical.target
